# deploy/discord-bot.service
[Unit]
Description=Discord Voice Mod Recorder Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/voice-mod-bot
ExecStart=/usr/bin/node /home/pi/voice-mod-bot/dist/index.js
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/home/pi/voice-mod-bot/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-bot

[Install]
WantedBy=multi-user.target
